# supervisor/conf.d/epp.conf
[program:epp]
directory=/opt/node-epp
command=/usr/local/bin/node-epp-rabbit
priority=150
redirect_stderr=true
stdout_logfile=syslog
